# tablet-mode.service
[Unit]
Description=Disable keyboard and touchpad for tablet mode

[Service]
ExecStart=/usr/bin/tablet-mode
StandardOutput=null
